# deploy/minio/crontab
# MinIO Backup Cron Schedule
# Run backup script daily at midnight
0 0 * * * /backup.sh
